# packages/telegram-remote/examples/launchd/com.example.gajae-code.telegram-remote.rpc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.example.gajae-code.telegram-remote.rpc</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/sh</string>
		<string>/Users/YOU/src/gajae-code/packages/telegram-remote/examples/launchd/bin/telegram-remote-rpc.sh</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>GJC_TELEGRAM_REMOTE_ENV_FILE</key>
		<string>/Users/YOU/Library/Application Support/gjc/telegram-remote.env</string>
		<key>GJC_TELEGRAM_REMOTE_CHECKOUT</key>
		<string>/Users/YOU/src/gajae-code</string>
		<key>GJC_BUN</key>
		<string>/Users/YOU/.bun/bin/bun</string>
		<key>GJC_TELEGRAM_REMOTE_SOCKET_WAIT_SEC</key>
		<string>30</string>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/YOU/src/gajae-code/packages/telegram-remote</string>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/YOU/Library/Logs/gjc-telegram-remote-rpc.out.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/YOU/Library/Logs/gjc-telegram-remote-rpc.err.log</string>
</dict>
</plist>
